# http-download.service.in
[Unit]
Description=HTTP/HTTPS Download Service with CURL

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/http-download
