# scp-transfer.service.in
[Unit]
Description=SCP Transfer Service

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/scp-transfer