# nvidia-secure-update/nvidia_ap_fw_updater.service.in
[Unit]
Description=Nvidia ap secure fw update service
After=multi-user.target

[Service]
Restart=always
ExecStart=/usr/bin/nvidia-fw-updater
SyslogIdentifier=nvidia-fw-updater

[Install]
WantedBy=multi-user.target
